# config/systemd/user/nk-scdaemon-reload.service
[Unit]
Description=Reload scdaemon for Nitrokey pin reset

[Service]
ExecStart=/usr/bin/gpgconf --reload scdaemon

[Install]
WantedBy=default.target
